# lib/systemd/system/upower.service
[Unit]
Description=Daemon for power management

[Service]
Type=dbus
BusName=org.freedesktop.UPower
ExecStart=/usr/lib/upower/upowerd
