# var/lib/dpkg/info/cron.conffiles
/etc/cron.hourly/.placeholder
/etc/cron.d/.placeholder
/etc/cron.daily/.placeholder
/etc/cron.daily/standard
/etc/init/cron.conf
/etc/cron.monthly/.placeholder
/etc/default/cron
/etc/crontab
/etc/pam.d/cron
/etc/cron.weekly/.placeholder
