# var/lib/dpkg/info/dhelp.list
/etc
/etc/apache2
/etc/apache2/conf.d
/etc/apache2/conf.d/dhelp.conf
/etc/dhelp.conf
/etc/cron.weekly
/etc/cron.weekly/dhelp
/etc/cron.monthly
/etc/cron.monthly/dhelp
