# var/lib/dpkg/info/dhelp.list
/.
/usr
/usr/share
/usr/share/dhelp
/usr/share/dhelp/scripts
/usr/share/dhelp/scripts/conv-pdftotext
/usr/share/dhelp/scripts/conv-pstotext
/usr/share/dhelp/scripts/gz-extract
/usr/share/dhelp/fetcher_error.rhtml
/usr/share/dhelp/index.rhtml
/usr/share/dhelp/section.rhtml
/usr/share/dhelp/search_error.tmpl
/usr/share/dhelp/search_pagination.tmpl
/usr/share/dhelp/search_results.tmpl
/usr/share/dhelp/swish++.conf
/usr/share/dhelp/index-deferred
/usr/share/doc
/usr/share/doc/dhelp
/usr/share/doc/dhelp/debian.jpg
/usr/share/doc/dhelp/folder.png
/usr/share/doc/dhelp/swirl.jpg
/usr/share/doc/dhelp/text.png
/usr/share/doc/dhelp/dhelp.css
/usr/share/doc/dhelp/THANKS
/usr/share/doc/dhelp/README.Debian
/usr/share/doc/dhelp/copyright
/usr/share/doc/dhelp/examples
/usr/share/doc/dhelp/examples/index_package_doc.rb
/usr/share/doc/dhelp/changelog.gz
/usr/share/doc/dhelp/NEWS.Debian.gz
/usr/share/locale
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/dhelp.mo
/usr/share/locale/el
/usr/share/locale/el/LC_MESSAGES
/usr/share/locale/el/LC_MESSAGES/dhelp.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/dhelp.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/dhelp.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/dhelp.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/dhelp.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/dhelp.mo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dhelp.1.gz
/usr/share/man/man8
/usr/share/man/man8/dhelp_parse.8.gz
/usr/share/menu
/usr/share/menu/dhelp
/usr/sbin
/usr/sbin/dhelp_parse
/usr/bin
/usr/bin/dhelp
/usr/lib
/usr/lib/cgi-bin
/usr/lib/cgi-bin/dsearch
/usr/lib/cgi-bin/dhelp_fetcher
/usr/lib/ruby
/usr/lib/ruby/1.8
/usr/lib/ruby/1.8/dhelp
/usr/lib/ruby/1.8/dhelp/exporter
/usr/lib/ruby/1.8/dhelp/exporter/base.rb
/usr/lib/ruby/1.8/dhelp/exporter/html.rb
/usr/lib/ruby/1.8/dhelp.rb
/var
/var/lib
/var/lib/dhelp
/var/lib/dhelp/tmp
/etc
/etc/apache2
/etc/apache2/conf.d
/etc/apache2/conf.d/dhelp.conf
/etc/apt
/etc/apt/apt.conf.d
/etc/dhelp.conf
/etc/cron.weekly
/etc/cron.weekly/dhelp
/etc/cron.monthly
/etc/cron.monthly/dhelp
