# var/lib/dpkg/info/dhelp.conffiles
/etc/apache2/conf.d/dhelp.conf
/etc/dhelp.conf
/etc/cron.weekly/dhelp
/etc/cron.monthly/dhelp
