# var/lib/dpkg/info/dlocate.list
/.
/usr
/usr/bin
/usr/bin/dlocate
/usr/sbin
/usr/sbin/update-dlocatedb
/usr/sbin/dpkg-hold
/usr/sbin/dpkg-unhold
/usr/sbin/dpkg-remove
/usr/sbin/dpkg-purge
/usr/share
/usr/share/doc
/usr/share/doc/dlocate
/usr/share/doc/dlocate/README.Debian
/usr/share/doc/dlocate/copyright
/usr/share/doc/dlocate/changelog.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dlocate.1.gz
/usr/share/man/man8
/usr/share/man/man8/dpkg-hold.8.gz
/usr/share/man/man8/dpkg-unhold.8.gz
/usr/share/man/man8/dpkg-remove.8.gz
/usr/share/man/man8/dpkg-purge.8.gz
/usr/share/man/man8/update-dlocatedb.8.gz
/var
/var/lib
/var/lib/dlocate
/etc
/etc/bash_completion.d
/etc/bash_completion.d/dlocate-completion
/etc/default
/etc/default/dlocate
/etc/cron.daily
/etc/cron.daily/dlocate
