# var/lib/dpkg/info/dlocate.conffiles
/etc/bash_completion.d/dlocate-completion
/etc/default/dlocate
/etc/cron.daily/dlocate
