# var/lib/dpkg/info/dwww.conffiles
/etc/dwww/apache.conf
/etc/cron.daily/dwww
/etc/cron.weekly/dwww
