# var/lib/dpkg/info/dwww.list
/.
/etc
/etc/dwww
/etc/dwww/apache.conf
/etc/cron.daily
/etc/cron.daily/dwww
/etc/cron.weekly
/etc/cron.weekly/dwww
/var
/var/cache
/var/cache/dwww
/var/cache/dwww/db
/var/lib
/var/lib/dwww
/var/www
/usr
/usr/sbin
/usr/sbin/dwww-txt2html
/usr/sbin/dwww-build
/usr/sbin/dwww-find
/usr/sbin/dwww-format-man
/usr/sbin/dwww-quickfind
/usr/sbin/dwww-convert
/usr/sbin/dwww-build-menu
/usr/sbin/dwww-refresh-cache
/usr/sbin/dwww-index++
/usr/sbin/dwww-cache
/usr/lib
/usr/lib/cgi-bin
/usr/lib/cgi-bin/dwww
/usr/bin
/usr/bin/dwww
/usr/share
/usr/share/iceape
/usr/share/iceape/searchplugins
/usr/share/applications
/usr/share/applications/dwww.desktop
/usr/share/iceweasel
/usr/share/iceweasel/searchplugins
/usr/share/dwww
/usr/share/dwww/templates
/usr/share/dwww/templates/man-index.end
/usr/share/dwww/templates/menu.errors.end
/usr/share/dwww/templates/man-in-section.start
/usr/share/dwww/templates/dwww.end
/usr/share/dwww/templates/man-begins-with.start
/usr/share/dwww/templates/menu.nosections
/usr/share/dwww/templates/dwww-convert.start
/usr/share/dwww/templates/menu.all.start
/usr/share/dwww/templates/man-begins-with.end
/usr/share/dwww/templates/menu.errors.start
/usr/share/dwww/templates/man-index.start
/usr/share/dwww/templates/dwww-find.start
/usr/share/dwww/templates/dwww-convert.dir.start
/usr/share/dwww/templates/menu.single.start
/usr/share/dwww/templates/dwww.template
/usr/share/dwww/templates/man-in-section.end
/usr/share/dwww/templates/menu.end
/usr/share/dwww/swish++.conf
/usr/share/dwww/mozilla
/usr/share/dwww/mozilla/dwww.png
/usr/share/dwww/mozilla/dwww.src
/usr/share/dwww/dwww.css
/usr/share/dwww/debian.png
/usr/share/dwww/functions.sh
/usr/share/menu
/usr/share/menu/dwww
/usr/share/doc
/usr/share/doc/dwww
/usr/share/doc/dwww/copyright
/usr/share/doc/dwww/changelog.gz
/usr/share/doc/dwww/README
/usr/share/doc/dwww/NEWS.Debian.gz
/usr/share/doc/dwww/TODO
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/dwww.7.gz
/usr/share/man/man8
/usr/share/man/man8/dwww-refresh-cache.8.gz
/usr/share/man/man8/dwww-format-man.8.gz
/usr/share/man/man8/dwww-find.8.gz
/usr/share/man/man8/dwww-cache.8.gz
/usr/share/man/man8/dwww-index++.8.gz
/usr/share/man/man8/dwww-txt2html.8.gz
/usr/share/man/man8/dwww-build.8.gz
/usr/share/man/man8/dwww-quickfind.8.gz
/usr/share/man/man8/dwww-convert.8.gz
/usr/share/man/man8/dwww-build-menu.8.gz
/usr/share/man/man1
/usr/share/man/man1/dwww.1.gz
/usr/share/man/man5
/usr/share/perl5
/usr/share/perl5/Debian
/usr/share/perl5/Debian/Dwww
/usr/share/perl5/Debian/Dwww/Version.pm
/usr/share/perl5/Debian/Dwww/ConfigFile.pm
/usr/share/perl5/Debian/Dwww/Common.pm
/usr/share/perl5/Debian/Dwww/DocBase.pm
/usr/share/perl5/Debian/Dwww/Utils.pm
/usr/share/perl5/Debian/Dwww/Initialize.pm
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/dwww
/usr/share/iceape/searchplugins/dwww.png
/usr/share/iceape/searchplugins/dwww.src
/usr/share/iceweasel/searchplugins/dwww.png
/usr/share/iceweasel/searchplugins/dwww.src
/usr/share/dwww/templates/dwww-convert.dir.end
/usr/share/dwww/templates/dwww-find.end
/usr/share/dwww/templates/dwww-convert.end
/usr/share/man/man5/dwww.conf.5.gz
