# var/lib/dpkg/info/fam.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/famd.8.gz
/usr/share/man/man5
/usr/share/man/man5/fam.conf.5.gz
/usr/share/doc
/usr/share/doc/fam
/usr/share/doc/fam/NEWS.Debian.gz
/usr/share/doc/fam/NEWS.gz
/usr/share/doc/fam/copyright
/usr/share/doc/fam/README
/usr/share/doc/fam/TODO
/usr/share/doc/fam/AUTHORS
/usr/share/doc/fam/changelog.Debian.gz
/usr/sbin
/usr/sbin/famd
/etc
/etc/init.d
/etc/init.d/fam
/etc/fam.conf
