# var/lib/dpkg/info/hostname.list
/.
/bin
/bin/domainname
/bin/dnsdomainname
/bin/nisdomainname
/bin/ypdomainname
/bin/hostname
/usr
/usr/share
/usr/share/man
/usr/share/man/fr
/usr/share/man/fr/man1
/usr/share/man/fr/man1/hostname.1.gz
/usr/share/man/man1
/usr/share/man/man1/hostname.1.gz
/usr/share/doc
/usr/share/doc/hostname
/usr/share/doc/hostname/changelog.gz
/usr/share/doc/hostname/copyright
/etc
/etc/init.d
/etc/init
/etc/init/hostname.conf
/usr/share/man/man1/dnsdomainname.1.gz
/usr/share/man/man1/nisdomainname.1.gz
/usr/share/man/man1/ypdomainname.1.gz
/usr/share/man/man1/domainname.1.gz
/etc/init.d/hostname
