# var/lib/dpkg/info/kexec-tools.conffiles
/etc/init.d/kexec
/etc/init.d/kexec-load
/etc/init.d/kdump
