# var/lib/dpkg/info/kexec-tools.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/kexec-tools
/usr/share/doc/kexec-tools/News
/usr/share/doc/kexec-tools/changelog.Debian.gz
/usr/share/doc/kexec-tools/TODO
/usr/share/doc/kexec-tools/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/kexec-tools
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/kexec.8.gz
/usr/share/man/man8/kdump.8.gz
/usr/share/man/man8/coldreboot.8.gz
/sbin
/sbin/kexec
/sbin/coldreboot
/sbin/kdump
/etc
/etc/init.d
/etc/init.d/kexec
/etc/init.d/kexec-load
/etc/init.d/kdump
