# var/lib/dpkg/info/lpr.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/lpr
/usr/share/doc/lpr/copyright
/usr/share/doc/lpr/README.Debian
/usr/share/doc/lpr/examples
/usr/share/doc/lpr/examples/printcap
/usr/share/doc/lpr/changelog.gz
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/printcap.5.gz
/usr/share/man/man8
/usr/share/man/man8/lpc.8.gz
/usr/share/man/man8/pac.8.gz
/usr/share/man/man8/lpd.8.gz
/usr/share/man/man8/lpf.8.gz
/usr/share/man/man1
/usr/share/man/man1/lpr.1.gz
/usr/share/man/man1/lpq.1.gz
/usr/share/man/man1/lptest.1.gz
/usr/share/man/man1/lprm.1.gz
/usr/sbin
/usr/sbin/lpf
/usr/sbin/lpd
/usr/sbin/pac
/usr/sbin/lpc
/usr/bin
/usr/bin/lprm
/usr/bin/lptest
/usr/bin/lpq
/usr/bin/lpr
/var
/var/log
/var/spool
/var/spool/lpd
/var/spool/lpd/lp
/var/spool/lpd/remote
/etc
/etc/default
/etc/default/lpd
/etc/cron.weekly
/etc/cron.weekly/lpr
/etc/init.d
/etc/init.d/lpd
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/lpr
