# var/lib/dpkg/info/man2html.list
/.
/var
/var/cache
/var/cache/man2html
/etc
/etc/cron.weekly
/etc/cron.weekly/man2html
/usr
/usr/share
/usr/share/doc-base
/usr/share/doc-base/man2html
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/hman.1.gz
/usr/share/man/el
/usr/share/man/el/man1
/usr/share/man/el/man1/hman.1.gz
/usr/share/man/it
/usr/share/man/it/man1
/usr/share/man/it/man1/hman.1.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/man2html
/usr/share/doc
/usr/share/doc/man2html-base
/usr/share/doc/man2html-base/html
/usr/share/doc/man2html-base/html/man.html
/usr/share/menu
/usr/share/menu/man2html
/usr/share/man2html
/usr/share/man2html/man.aux
/usr/share/man2html/swish++.conf
/usr/share/man2html/mansearch.aux
/usr/bin
/usr/bin/hman
/usr/lib
/usr/lib/cgi-bin
/usr/lib/cgi-bin/man
/usr/lib/cgi-bin/man/manwhatis
/usr/lib/cgi-bin/man/mansec
/usr/lib/cgi-bin/man/mansearch
/usr/lib/cgi-bin/man/man2html
/usr/share/doc/man2html
