# var/lib/dpkg/info/mercurial.conffiles
/etc/bash_completion.d/mercurial
/etc/mercurial/hgrc
/etc/mercurial/hgrc.d/mergetools.rc
/etc/mercurial/hgrc.d/cacerts.rc
