# var/lib/dpkg/info/mercurial.list
/.
/usr
/usr/lib
/usr/lib/python2.7
/usr/lib/python2.7/dist-packages
/usr/lib/python2.7/dist-packages/mercurial
/usr/lib/python2.7/dist-packages/mercurial/mpatch.so
/usr/lib/python2.7/dist-packages/mercurial/parsers.so
/usr/lib/python2.7/dist-packages/mercurial/osutil.so
/usr/lib/python2.7/dist-packages/mercurial/base85.so
/usr/lib/python2.7/dist-packages/mercurial/diffhelpers.so
/usr/lib/python2.7/dist-packages/mercurial/bdiff.so
/usr/lib/python2.7/dist-packages/hgext
/usr/lib/python2.7/dist-packages/hgext/inotify
/usr/lib/python2.7/dist-packages/hgext/inotify/linux
/usr/lib/python2.7/dist-packages/hgext/inotify/linux/_inotify.so
/usr/lib/pyshared
/usr/lib/pyshared/python2.7
/usr/lib/pyshared/python2.7/mercurial
/usr/lib/pyshared/python2.7/hgext
/usr/lib/pyshared/python2.7/hgext/inotify
/usr/lib/pyshared/python2.7/hgext/inotify/linux
/usr/share
/usr/share/mercurial
/usr/share/mercurial/hgext.rc.md5sums
/usr/share/mercurial/hgext.rc
/usr/share/doc
/usr/share/doc/mercurial
/usr/share/doc/mercurial/changelog.Debian.gz
/usr/share/doc/mercurial/NEWS.Debian.gz
/usr/share/doc/mercurial/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mercurial
/usr/bin
/usr/bin/hg
/etc
/etc/bash_completion.d
/etc/bash_completion.d/mercurial
/etc/mercurial
/etc/mercurial/hgrc
/etc/mercurial/hgrc.d
/etc/mercurial/hgrc.d/mergetools.rc
/etc/mercurial/hgrc.d/cacerts.rc
/usr/lib/pyshared/python2.7/mercurial/mpatch.so
/usr/lib/pyshared/python2.7/mercurial/parsers.so
/usr/lib/pyshared/python2.7/mercurial/osutil.so
/usr/lib/pyshared/python2.7/mercurial/base85.so
/usr/lib/pyshared/python2.7/mercurial/diffhelpers.so
/usr/lib/pyshared/python2.7/mercurial/bdiff.so
/usr/lib/pyshared/python2.7/hgext/inotify/linux/_inotify.so
/usr/share/doc/mercurial/examples
