# var/lib/dpkg/info/setserial.conffiles
/etc/modutils/setserial
/etc/init.d/setserial
/etc/init.d/etc-setserial
